# systemd/ffl-auto-sync.service
[Unit]
Description=Firefighter Linux Auto Sync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ffl-auto-sync
